# deploy/systemd/jicho-ui.service
[Unit]
Description=Jicho Smart UI Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tinybeaver
WorkingDirectory=/home/tinybeaver/traffic_system
Environment=PYTHON_BIN=/usr/bin/python3
Environment=DASHBOARD_HOST=0.0.0.0
Environment=DASHBOARD_PORT=3000
EnvironmentFile=/home/tinybeaver/traffic_system/.env
ExecStart=${PYTHON_BIN} -m uvicorn dashboard:app --host ${DASHBOARD_HOST} --port ${DASHBOARD_PORT}
Restart=always
RestartSec=5
StandardOutput=append:/home/tinybeaver/traffic_system/logs/ui.log
StandardError=append:/home/tinybeaver/traffic_system/logs/ui.log

[Install]
WantedBy=multi-user.target
