# deploy/systemd/jicho-detector.service
[Unit]
Description=Jicho Smart Detector Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__REPO_DIR__
Environment=PYTHON_BIN=/usr/bin/python3
EnvironmentFile=__REPO_DIR__/.env
Environment=RUN_DASHBOARD_IN_MAIN=0
ExecStart=${PYTHON_BIN} __REPO_DIR__/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jicho-detector

[Install]
WantedBy=multi-user.target
